# devsecops-routing-basics-solution/router1/_system/etc/rc.local
#!/bin/bash
sudo ip route del default
sudo ip route add 10.1.2.0/24 via 172.16.0.2
