# sparta-nat/gateway/_system/etc/rc.local
#!/bin/bash

sudo ip route del default
sudo ip route add default via 200.200.200.1
