# devsecops-routing-basics-2/devsecops-routing-basics-2/router2/_system/etc/rc.local
#!/bin/bash

NETWORK2_IF=$(ifconfig | grep -B1 "inet addr:10.1.2.254" | awk '$1!="inet" && $1!="--" {print $1}')
CONNECTION_NET_IF=$(ifconfig | grep -B1 "inet addr:172.16.0.2" | awk '$1!="inet" && $1!="--" {print $1}')

ip add flush $NETWORK2_IF
ip add flush $CONNECTION_NET_IF
sudo ip route del default
